# docker/config/dev/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:consumer]
command=/var/www/html/yii kafka/consume
autostart=true
autorestart=true
